# example/pdns/supervisord.conf
[supervisord]
nodaemon=true

[program:syslog]
command=/sbin/syslogd

[program:pdns]
command=/usr/sbin/pdns_server --guardian=yes